# config/crontabs/root
# This is crontab file. Here you can customize your cron-jobs.
# Changes will apply after file save. Container restart not need
#

# * * * * * - time options, see more https://crontab.guru/
# flock - lock files are used to determine if a script or application is already running
# >> - append, > renew log file
# >>/var/log/cron.log 2>&1 - send output to log
# >/dev/null 2>&1 - send output to black hole
# 2>&1 redirects the output include both the Standard Error and Standard Out. Look details https://tldp.org/LDP/abs/html/io-redirection.html

# Example
# * * * * * [flock -n <path-to-lck-file>] <some-command> >><path-to-log-file> 2>&1
# * * * * * [flock -n <path-to-lck-file>] <some-command> >/dev/null 2>&1



# Let's run!

# WordPress cron run by system cron
*/1 * * * * docker exec ${APP_NAME}_wordpress flock -n /var/log/wordpress/wp-cron.lck php /var/www/html/wp-cron.php >> /var/log/cron.log 2>&1

# Remove wp-cron log file once a week
#1 3 * * 2 rm /var/log/cron/wp-cron.log >/dev/null 2>&1

# An empty line is required at the end of this file for a valid cron file