# config/crontabs/root
*/1 * * * * echo 'Happy New Year!!' >> /var/log/cron.log 2>&1